# filesystem/etc/inittab
::sysinit:/etc/init.d/rc
::askfirst:-/bin/sh /home/main
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
