# filesystem/etc/inittab
::sysinit:/etc/init.d/rc
::askfirst:/bin/sh /app/init
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount­-a -r
